# tools-setup/roles/sonarqube/templates/sonarqube.service
[Unit]
Description=SonarQube Server
After=network.target

[Service]
ExecStart=/sonarqube/sonarqube-25.2.0.102705/bin/linux-x86-64/sonar.sh start
User=sonar
KillMode=process
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
